# ops/systemd/sofascore-historical-hydrate@.service
[Unit]
Description=Sofascore Historical Hydrate Worker %i
After=network-online.target redis-server.service postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/sofascore
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/sofascore/deploy/run_service.sh worker-historical-hydrate --consumer-name historical-hydrate-%i --block-ms 5000
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=90

[Install]
WantedBy=multi-user.target
