# ops/systemd/sofascore-live-warm@.service
[Unit]
Description=Sofascore Live Warm Worker %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/sofascore
EnvironmentFile=/opt/sofascore/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/sofascore/.venv/bin/python -m schema_inspector.cli --log-level INFO worker-live-warm --consumer-name live-warm-%i --block-ms 5000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
