# ops/systemd/sofascore-historical-tournament-planner.service
[Unit]
Description=Sofascore Historical Tournament Planner Daemon
After=network-online.target redis-server.service postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/sofascore
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/sofascore/deploy/run_service.sh historical-tournament-planner-daemon
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
