# ops/systemd/sofascore-structure-sync.service
[Unit]
Description=Sofascore Structure Sync Worker
After=network-online.target redis-server.service postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/sofascore
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/sofascore/deploy/run_service.sh worker-structure-sync --consumer-name structure-sync-1 --block-ms 5000
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
