# ops/systemd/sofascore-discovery@.service
[Unit]
Description=Sofascore Discovery Worker %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/sofascore
EnvironmentFile=/opt/sofascore/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/sofascore/.venv/bin/python -m schema_inspector.cli --log-level INFO worker-discovery --consumer-name discovery-%i --block-ms 5000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
